# patches/etc/rc.local
#!/bin/sh

/usr/local/bin/install_mysql_server.sh

exit 0
